# conf.d/xvnc.conf
[program:xvnc]
command=Xvnc :0 -geometry %(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)s -depth 16
autorestart=true
